# bdgd.service
[Service]
ExecStart=/home/gridflow/src/venv/bin/python populate_db.py
WorkingDirectory=/home/gridflow/src